# templates/default/livy.service.erb
[Unit]
Description=Livy server
After=syslog.target network.target remote-fs.target <%= @deps %>
#restartSec is 5 
#a maximum burst of 121 in 11min guaranty that it will at least try for 10min 
StartLimitIntervalSec=660
StartLimitBurst=121

[Service]
User = <%= node['livy']['user'] %>

PIDFile=<%= node['livy']['pid_file'] %>
ExecStart=<%= node['livy']['base_dir'] %>/bin/start-livy.sh
ExecStop=<%= node['livy']['base_dir'] %>/bin/stop-livy.sh
LimitNOFILE=20000
TimeoutSec= 10min
Restart=on-failure
RestartSec=5

Type=forking

[Install]
WantedBy=multi-user.target
